# meta-xbterminal/recipes-xbterminal/xbterminal-firmware/files/xbterminal-x.service
[Unit]
Description=XBTerminal X Service
After=syslog.target
PartOf=xbterminal-firmware.service

[Service]
Type=simple
ExecStart=/usr/bin/X -nocursor
Nice=-3

KillMode=control-group
TimeoutStopSec=600

Restart=on-failure
RestartSec=30
StartLimitInterval=330
StartLimitBurst=5



[Install]
WantedBy=multi-user.target
